# autostart/ibus-daemon.desktop
[Desktop Entry]
Type=Application
Name=iBus Daemon
TryExec=ibus-daemon
Exec=ibus-daemon -drx
StartupNotify=false
