# autostart/kdeconnectd.desktop
[Desktop Entry]
Type=Application
Name=KDEConnect Daemon
TryExec=/usr/lib/kdeconnectd
Exec=/usr/lib/kdeconnectd
StartupNotify=false
