# autostart/mpd.desktop
[Desktop Entry]
Type=Application
Name=MPD
TryExec=mpd
Exec=mpd
StartupNotify=false
