# autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=Network Manager Applet
TryExec=nm-applet
Exec=nm-applet
StartupNotify=false
