# systemd/battery-monitor.timer
[Unit]
Description=Monitors battery level and charging status

[Timer]
OnActiveSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
